# backup/crontabs
0 0 * * * bash ~/enginetribe-misc/backup/perform-backup.sh
